# blog-store-messages.service.in
[Unit]
Description=Tell blogd about writable boot log
DefaultDependencies=no
RequiresMountsFor=@@BOOT_LOGFILE@@
Conflicts=shutdown.target
After=systemd-remount-fs.service auditd.service
Before=sysinit.target shutdown.target
ConditionPathExists=!/etc/initrd-release

[Service]
Type=oneshot
TimeoutSec=20
RemainAfterExit=yes
ExecStart=-/sbin/blogctl ready
ExecStop=-/sbin/blogctl close

[Install]
WantedBy=sysinit.target halt.target kexec.target poweroff.target reboot.target
