# systemd-ask-password-blog.path
[Unit]
Description=Forward Password Requests to Plymouth Directory Watch
Documentation=http://www.freedesktop.org/wiki/Software/systemd/PasswordAgents
DefaultDependencies=no
Conflicts=shutdown.target
After=blog.service
Before=basic.target shutdown.target
ConditionKernelCommandLine=!blog.enable=0
ConditionPathExists=/run/blogd.pid

[Path]
DirectoryNotEmpty=/run/systemd/ask-password
MakeDirectory=yes
